# debian/storage-server.cron.d
* * * * * root cd /var/lib/storage-server/ && /var/lib/storage-server/watch-dog.sh >> /var/lib/storage-server/storage-server.log 2>&1
